# scripts/deployment/com.alphastock.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphastock.scheduler</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/adithyasaladi/Personal/AlphaStock/.venv/bin/python</string>
        <string>/Users/adithyasaladi/Personal/AlphaStock/scheduler.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/adithyasaladi/Personal/AlphaStock</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/adithyasaladi/Personal/AlphaStock/logs/scheduler.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/adithyasaladi/Personal/AlphaStock/logs/scheduler_error.log</string>
</dict>
</plist>
